# containers/penpot/penpot-redis.container
[Container]
ContainerName=penpot-redis
Image=redis:7.2
Network=penpot.network

[Service]
Restart=always
